# .config/systemd/user/foot-server.socket
[Socket]
ListenStream=%t/foot-%i.sock

[Install]
WantedBy=default.target
