# .config/systemd/user/xremap.service
[Unit]
Description=xremap

[Service]
KillMode=process
ExecStart=/home/fuzakebito/.cargo/bin/xremap --watch /home/fuzakebito/.config/xremap/config.yml
Type=simple
Restart=always
RestartSec=2s

Environment=DISPLAY=:0

[Install]
WantedBy=graphical-session.target
